# docker/cron/periodic/15min/genova_docker_cleanup
#!/bin/sh

cd /app
bundle exec thor genova:docker-cleanup >> /app/log/cron.log 2>&1
