# docker/cron/crontabs/root
20 0 * * * /usr/sbin/logrotate /etc/logrotate.conf >> /app/log/cron.log 2>&1
10 0 * * * /usr/sbin/logrotate -f /etc/logrotate.d/genova
