# docker/cron/periodic/hourly/genova
#!/bin/sh

cd /app

echo 'thor genova:util:delete_invalid_jobs'
thor genova:utils delete-invalid-jobs

echo 'thor genova:util:delete_old_images'
thor genova:utils delete-old-images

