# docker/cron/periodic/hourly/genova
#!/bin/sh

cd /app
bundle exec thor genova:docker-cleanup
